# packages/mykitchen/mykitchen-0.0.1.tar.gz/mykitchen-0.0.1/deploy/etc/s6-overlay/s6-rc.d/worker/run
#!/command/with-contenv execlineb

cd /app
exec python mykitchen/__main__.py qcluster